# setup/statusscraper.cron
# Execute statusscraper every minute
* * * * *   root   /srv/www/statusscraper/statusscraper.sh